# home/Library/LaunchAgents/com.r3nniw.pi-antigravity-rotator.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.r3nniw.pi-antigravity-rotator</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/homebrew/sbin/sleepwatcher</string>
      <string>--verbose</string>

      <!-- system sleep → kill PAR -->
      <string>--sleep</string>
      <string>{{ .chezmoi.homeDir }}/.local/bin/par-manager _on-sleep</string>

      <!-- system wake → start PAR -->
      <string>--wakeup</string>
      <string>{{ .chezmoi.homeDir }}/.local/bin/par-manager _on-wake</string>

      <!-- display sleep → kill PAR -->
      <string>--displaysleep</string>
      <string>{{ .chezmoi.homeDir }}/.local/bin/par-manager _on-sleep</string>

      <!-- display wake → start PAR -->
      <string>--displaywakeup</string>
      <string>{{ .chezmoi.homeDir }}/.local/bin/par-manager _on-wake</string>

      <!-- idle timeout (150 = 15s in 0.1s units; 36000 = 1hr) -->
      <string>--timeout</string>
      <string>36000</string>

      <!-- after idle timeout -> kill PAR -->
      <string>--idle</string>
      <string>{{ .chezmoi.homeDir }}/.local/bin/par-manager _on-idle</string>

      <!-- resume from idle → start PAR -->
      <string>--idleresume</string>
      <string>{{ .chezmoi.homeDir }}/.local/bin/par-manager _on-wake</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
      <key>PATH</key>
      <string>{{ .chezmoi.homeDir }}/.local/share/mise/installs/node/latest/bin:{{ .chezmoi.homeDir }}/.local/share/mise/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/com.r3nniw.pi-antigravity-rotator.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.r3nniw.pi-antigravity-rotator.err.log</string>
  </dict>
</plist>
